# 99-miix.rules
ATTR{idVendor}=="17ef", ATTR{idProduct}=="60bb", RUN+="/usr/local/sbin/miix_fnfix.py"
